# osiris/mistral-small-4-vllm/mistral-small-4-vllm-watchdog.timer
[Unit]
Description=Check Mistral Small 4 vLLM health every minute

[Timer]
OnBootSec=20min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=mistral-small-4-vllm-watchdog.service

[Install]
WantedBy=timers.target
