# osiris/mistral-small-4-vllm/mistral-small-4-vllm.service
[Unit]
Description=Mistral Small 4 119B NVFP4 vLLM on Osiris
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
Conflicts=gemma4-26b-a4b-vllm.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/royossai/services/mistral-small-4-vllm
ExecStart=/usr/bin/docker compose -f /home/royossai/services/mistral-small-4-vllm/compose.yaml up -d
ExecStop=/usr/bin/docker compose -f /home/royossai/services/mistral-small-4-vllm/compose.yaml down
TimeoutStartSec=0
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target
