# osiris/gemma4-26b-a4b-vllm/gemma4-26b-a4b-vllm-watchdog.service
[Unit]
Description=Restart Gemma 4 26B A4B vLLM if the API is unhealthy
After=docker.service network-online.target gemma4-26b-a4b-vllm.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/royossai/services/gemma4-26b-a4b-vllm
ExecStart=/bin/sh -c 'curl -fsS http://127.0.0.1:8000/health >/dev/null || /usr/bin/docker compose -f /home/royossai/services/gemma4-26b-a4b-vllm/compose.yaml restart gemma4-26b-a4b-vllm'
